# systemd/auth.service
[Unit]
Description=SMTP Authentication Service

[Service]
ExecStart=/usr/local/sbin/auth
Restart=always
