# systemd/forwarding.service
[Unit]
Description=SMTP Forwarding Service

[Service]
ExecStart=/usr/local/sbin/forwarding
Restart=always
