# systemd/mailout.service
[Unit]
Description=Mail Outgoing Gateway

[Service]
ExecStart=/usr/local/sbin/mailout
Restart=always
