# systemd/maildb.service
[Unit]
Description=Mail Database Service

[Service]
ExecStart=/usr/local/sbin/maildb
Restart=always
