# systemd/mailway-supervisor.service
[Unit]
Description=Mailway supervisor

[Service]
ExecStart=/usr/local/sbin/mailway supervisor
Restart=always
