# systemd/webhooks.service
[Unit]
Description=Webhooks service

[Service]
ExecStart=/usr/local/sbin/webhooks
Restart=always
